# Server/redbox.cron
#!/bin/sh
sudo -u redbox /home/redbox/deploy.sh redbox
sudo -u redbox /home/redbox/deploy.sh mint